# deployment/cadre-package-listener.supervisor.ini
[program:cadre-job-listener]
command = /home/ubuntu/cadre-job-listener/deployment/package_listener_start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/cadre-package-listener.super.err.log
stdout_logfile=/var/log/cadre-package-listener.super.out.log
